# ansible/roles/dns2/files/cron.d/dns_root_agent
# m     h       dom     mon     dow     user    command
*       *       *       *       *       root    /usr/local/bin/dns_root_agent_cronjob
