# ansible/roles/spy/files/cron.d/update_debian
# m     h       dom     mon     dow     user    command
17      13      *       *       *       root    /usr/local/bin/update_debian_cronjob
