# ansible/roles/print_server/files/bin/pdfprint
#!/bin/bash

### BEGIN INIT INFO
# Provides:          pdfprint
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start pdfprint at boot time
# Description:       Print-to-PDF voodoo to present a consistent virtual printer front-end for a physical printer back-end.
### END INIT INFO


function ssd {
  start-stop-daemon \
    -c rachel \
    -u rachel \
    -a /usr/local/bin/pdfprintd \
    -n pdfprintd \
    "${@}"
}

case "${1}" in
  'status')
    ssd -T
    ;;
  'stop')
    ssd -K
    ;;
  'start') 
    ssd -b -S
    ;;
  'restart')
    ssd -K
    ssd -b -S
    ;;
  '*')
    echo "Valid sub-commands are:  status stop start restart"
    exit 1
    ;;
esac

