# cron.d/meliora_reboot
# m     h       dom     mon     dow     user    command
45      00      *       *       *       root    /usr/sbin/shutdown -r +5 "scheduled reboot"
