# cron.d/reboot
# Reboot every day at noon.  Kick in updates and what not.

# m     h       dom     mon     dow     user    command
45      23      *       *       *       root    /usr/sbin/shutdown -r +1 "scheduled reboot"
