# cron.d/update_live-build_meliora
# m     h       dom     mon     dow     user    command
15      01      *       *       *       root    /usr/local/bin/update_live-build_cronjob
