# playbooks/roles/dns2/files/cron.d/update_dns
# m     h       dom     mon     dow     user    command
*       *       *       *       *       pi      /usr/local/bin/dns_root_agent
