# ansible/roles/backup_client/templates/cron.d/backup.j2
# m     h       dom     mon     dow     user    command
{% for item in backup_users %}
00      *       *       *       *       {{ item }} /usr/local/bin/backup_cronjob
{% endfor %}
