# ansible/roles/storage/files/cron.d/tend_storage
# m     h       dom     mon     dow     user    command
@reboot                                 root    /usr/local/bin/tend_storage > /var/log/tend_storage.log 2>&1
00      *       *       *       *       root    /usr/local/bin/tend_storage > /var/log/tend_storage.log 2>&1
