# ansible/roles/dst/files/cron.d/reboot
# m     h       dom     mon     dow     user    command
57      10      *       *       *       root    /usr/sbin/shutdown -r +1 "scheduled reboot"
