# ansible/roles/dst/files/cron.d/update_debian
# m     h       dom     mon     dow     user    command
32      02      *       *       *       root    /usr/local/bin/update_debian_cronjob
