# ansible/roles/meliora/files/cron.d/reboot
# m     h       dom     mon     dow     user    command
29      01      *       *       *       root    /sbin/shutdown -r +1 "scheduled reboot"
